# nbd-runner.service.in
[Unit]
Description=Network Block Device For Distributed Storages
Requires=rpcbind.service
Wants=network-online.target
After=network-online.target rpcbind.service

[Service]
Type=simple
Environment="NOFILE=1048576"
EnvironmentFile=-/etc/sysconfig/nbd-runner
ExecStart=@prefix@/sbin/nbd-runner --threads 2
ExecStartPost=-/bin/bash -c "prlimit --pid $MAINPID --core=unlimited --nofile=$NOFILE:$NOFILE"
KillMode=process
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
